# contrib/systemd/clip-sync-client-sync.service
[Unit]
Description=Clip-Sync Client Sync (user)
After=default.target

[Service]
Type=simple
WorkingDirectory=%h/Work/clip-sync
# Ajusta token/device si quieres
ExecStart=%h/Work/clip-sync/bin/cli --mode sync --addr ws://127.0.0.1:8080/ws --token u1 --device L1 --poll-ms 500
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target

